# deploy/systemd/pi-audit-backup.service
# SUPERSEDED (issue #201-followup) for THIS repo's actual (Docker-Compose-based) demo
# deployment — the demo now runs the `audit-backup` compose service in
# docker/docker-compose.demo.yml instead (Docker-native: no host bind-mount, no host root
# needed to install a timer, scheduled internally via busybox crond — see
# docker/backup/audit-backup/ and deploy/README.md). This unit is kept, UNMODIFIED otherwise,
# because it documents a real prior design and may still be relevant for a non-Compose
# deployment style (a bare-host Postgres reached via `docker exec`). Do not install this on
# the hp04 demo host — it would race the compose service for the same PI_BACKUP_DIR.
#
# Nightly logical backup of the BFF's audit Postgres (IMPROVEMENT-PLAN P0 #4 / Q4).
# Install on the demo host (hp04):
#   sudo cp deploy/systemd/pi-audit-backup.* /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now pi-audit-backup.timer
# Override the destination (a SECOND disk) + retention via a drop-in:
#   sudo systemctl edit pi-audit-backup.service   # [Service] Environment=PI_BACKUP_DIR=/mnt/backups/pi-audit
[Unit]
Description=Process Inspector — nightly audit-DB logical backup
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
Environment=PI_BACKUP_DIR=/var/backups/pi-audit
# Adjust WorkingDirectory to the repo checkout on the host.
WorkingDirectory=/home/flapci/workspace/process-inspector
ExecStart=/home/flapci/workspace/process-inspector/deploy/backup-audit-db.sh
