# deploy/systemd/pi-audit-backup.timer
[Unit]
Description=Run the Process Inspector audit-DB backup nightly

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true

[Install]
WantedBy=timers.target
